# install/ap-guardian.service
[Unit]
Description=AP-Guardian - Active Network Protection System
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
# Используем wrapper скрипт (создается автоматически или через setup.py)
ExecStart=/usr/local/bin/ap-guardian
# Альтернатива: ExecStart=/usr/bin/python3 -m src.main
# Рабочая директория
WorkingDirectory=/usr/lib/ap-guardian
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Ограничения ресурсов
MemoryLimit=100M
CPUQuota=30%

# Переменные окружения
Environment="PYTHONPATH=/usr/lib/ap-guardian"
Environment="AP_GUARDIAN_INTERFACE=any"

[Install]
WantedBy=multi-user.target
